# sppix-django.service
# SPPIX Django Systemd Service
# This service runs the Django application

[Unit]
Description=SPPIX Django Application
Documentation=https://docs.djangoproject.com/
After=network.target mysql.service redis.service
Wants=mysql.service redis.service

[Service]
Type=exec
User=sppix
Group=sppix
WorkingDirectory=/opt/sppix-store/Backend
Environment=PATH=/opt/sppix-store/Backend/venv/bin
Environment=DJANGO_SETTINGS_MODULE=core.settings
Environment=PYTHONPATH=/opt/sppix-store/Backend
Environment=DJANGO_ENV=production

# Main Django server
ExecStart=/opt/sppix-store/Backend/venv/bin/python /opt/sppix-store/Backend/manage.py runserver 127.0.0.1:82

# Restart policy
Restart=always
RestartSec=3
StartLimitInterval=60
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/sppix-store/Backend
ReadWritePaths=/opt/sppix-store/Backend/logs
ReadWritePaths=/opt/sppix-store/Backend/media
ReadWritePaths=/opt/sppix-store/Backend/staticfiles

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sppix-django

[Install]
WantedBy=multi-user.target
