# sppix-asgi.service
# SPPIX Django ASGI Systemd Service
# This service runs the Django ASGI application for WebSocket support

[Unit]
Description=SPPIX Django ASGI Application
Documentation=https://channels.readthedocs.io/
After=network.target mysql.service redis.service
Wants=mysql.service redis.service

[Service]
Type=exec
User=sppix
Group=sppix
WorkingDirectory=/opt/sppix-store/Backend
Environment=PATH=/opt/sppix-store/Backend/venv/bin
Environment=DJANGO_SETTINGS_MODULE=core.settings
Environment=PYTHONPATH=/opt/sppix-store/Backend
Environment=DJANGO_ENV=production

# ASGI server for WebSocket support
ExecStart=/opt/sppix-store/Backend/venv/bin/python /opt/sppix-store/Backend/run_asgi_server.py

# Restart policy
Restart=always
RestartSec=3
StartLimitInterval=60
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/sppix-store/Backend
ReadWritePaths=/opt/sppix-store/Backend/logs
ReadWritePaths=/opt/sppix-store/Backend/media
ReadWritePaths=/opt/sppix-store/Backend/staticfiles

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sppix-asgi

[Install]
WantedBy=multi-user.target
